# thief_server/thief-server.service
[Unit]
Description=Thief Game Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/thief_server
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
